# ansible/roles/docker/install/files/docker.service.d/exec.conf
[Service]
EnvironmentFile=-/etc/sysconfig/docker
ExecStart=
ExecStart=/usr/bin/docker daemon --storage-driver=overlay -H fd:// \
          $INSECURE_REGISTRY
