# ansible/roles/grafana/ansible-webhook/templates/grafana-ansible-webhook.socket.j2
[Unit]
Description=Grafana Ansible Webhook

[Socket]
ListenStream=127.0.0.1:8888
Accept=yes

[Install]
WantedBy=sockets.target
